# fix-eth-name.service
[Unit]
Description=Fix Eth Name
# Ensure the script exists before trying to run
ConditionPathExists=/fix-eth-name.sh

# Ensure udev has settled device names and local filesystems are up
After=systemd-udev-settle.service local-fs.target
Wants=systemd-udev-settle.service

# Run before various network services try to use potentially wrong interface names
# Added network-pre.target which is an early target before network configuration.
Before=network-pre.target network.service networking.service systemd-networkd.service NetworkManager.service wicked.service connman.service dhcpcd.service cloud-init-local.service
Before=network.target

[Service]
Type=oneshot
# Use absolute paths for commands for robustness
ExecStart=/usr/bin/env bash /fix-eth-name.sh
# Cleanup commands
ExecStart=/bin/rm -f /fix-eth-name.sh
ExecStart=/bin/rm -f /etc/systemd/system/fix-eth-name.service
ExecStart=/bin/rm -f /etc/systemd/system/multi-user.target.wants/fix-eth-name.service
# Check both potential preset directories
ExecStart=/bin/rm -f /lib/systemd/system-preset/01-fix-eth-name.preset
ExecStart=/bin/rm -f /usr/lib/systemd/system-preset/01-fix-eth-name.preset
# RemainAfterExit=yes ensures cleanup happens even if script fails, but might mask issues.
# StandardOutput=journal ensures script output goes to journald.
StandardOutput=journal+console
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
